# package/debian/DEBIAN/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]
then
  groupadd -r thingsplex || true
  groupadd -r easee || true

  useradd -r -N -M -d / -s /bin/false -g thingsplex easee || true

  chown -R easee:thingsplex /var/log/thingsplex/easee/ || true
  chown -R easee:easee /opt/thingsplex/easee/ || true
fi

deb-systemd-helper unmask easee.service >/dev/null || true

if deb-systemd-helper --quiet was-enabled easee.service; then
  systemctl enable easee.service || true
else
  deb-systemd-helper update-state easee.service >/dev/null || true
fi

if [ -d /run/systemd/system ]; then
  systemctl --system daemon-reload >/dev/null || true

  if [ -n "$2" ]; then
    deb-systemd-invoke try-restart easee.service >/dev/null || true
  else
    deb-systemd-invoke start easee.service >/dev/null || true
  fi
fi